# systemd/downdetector-refresh.timer
[Unit]
Description=Refresh Downdetector cache every 5 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min
AccuracySec=30s
Persistent=true
Unit=downdetector-refresh.service

[Install]
WantedBy=timers.target
