# systemd/downdetector-refresh.service
[Unit]
Description=Refresh Downdetector BR cache (parallel FlareSolverr)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
User=zabbix
Group=zabbix
Nice=10
EnvironmentFile=-/etc/zabbix/downdetector-api.env
# 4 sessões FlareSolverr em paralelo (ajuste DOWNDETECTOR_WORKERS se quiser)
Environment=DOWNDETECTOR_WORKERS=4
ExecStart=/usr/bin/python3 /opt/downdetector-zabbix/downdetector_scraper.py --refresh-all --workers 4 --delay 0
TimeoutStartSec=7200
