# systemd/downdetector-refresh.service
[Unit]
Description=Refresh Downdetector BR cache (serial FlareSolverr)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
User=zabbix
Group=zabbix
Nice=10
EnvironmentFile=-/etc/zabbix/downdetector-api.env
# Lista de hosts = API Zabbix (template downdetector). Agent só lê cache.
ExecStart=/usr/bin/python3 /opt/downdetector-zabbix/downdetector_scraper.py --refresh-all --delay 1.5
TimeoutStartSec=7200
