# systemd/downdetector-refresh.timer
[Unit]
Description=Refresh Downdetector BR cache every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=1min
Persistent=true
Unit=downdetector-refresh.service

[Install]
WantedBy=timers.target
