# systemd/downdetector-refresh.service
[Unit]
Description=Refresh Downdetector BR cache via FlareSolverr
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
User=zabbix
Group=zabbix
Nice=10
EnvironmentFile=-/etc/zabbix/downdetector-api.env
# 2 workers = bom equilíbrio velocidade/CPU (suba para 3–4 se a máquina aguentar)
Environment=DOWNDETECTOR_WORKERS=2
ExecStart=/usr/bin/python3 /opt/downdetector-zabbix/downdetector_scraper.py --refresh-all --workers 2
TimeoutStartSec=7200
